# templates/stardog.service.erb
<%-
ssl = scope['stardog::ssl']
scheme = (ssl == 'require' or ssl == 'enable') ? 'snarls' : 'snarl'
start_opt= case scope['stardog::ssl']
  when 'require' then '--require-ssl'
  when 'enable'  then '--enable-ssl' 
  else ''
end
-%>
# This file is managed by the Puppet <%= @name %> module

[Unit]
Description=The Stardog Server
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
EnvironmentFile=<%= scope['stardog::sysconfig'] %>
User=<%= scope['stardog::user'] %>
ExecStart=<%= scope['stardog::base'] %>/bin/stardog-admin server start <%= start_opt %>
ExecStop=<%= scope['stardog::base'] %>/bin/stardog-admin --server <%= scheme %>://127.0.0.1:<%= scope['stardog::port'] %>/ server stop
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
