# templates/stardog.service.erb
# This file is managed by the Puppet <%= @name %> module

[Unit]
Description=The Stardog Server
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
EnvironmentFile=<%= scope['stardog::sysconfig'] %>
User=<%= scope['stardog::user'] %>
ExecStart=<%= scope['stardog::base'] %>/bin/stardog-admin server start
ExecStop=<%= scope['stardog::base'] %>/bin/stardog-admin server stop
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
